# anpr.service
[Unit]
Description=ANPR API Service
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/anpr-system/anpr-system-main
Environment="PATH=/home/ubuntu/anpr-system/venv/bin"
ExecStart=/home/ubuntu/anpr-system/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
